# roles/openslides-add-instance/templates/openslides_instance.socket.j2
[Socket]
ListenStream={{ openslides_instance_port }}

[Install]
WantedBy=sockets.target
